# windows_linux/packaging/linux/deb/DEBIAN/postinst
#!/usr/bin/env bash
set -eu

remove_legacy_local_gui_wrapper() {
  local local_wrapper="/usr/local/bin/glasstokey-gui"
  if [ ! -f "${local_wrapper}" ]; then
    return
  fi

  if grep -Fq 'exec "/opt/GlassToKey.Linux.Gui/GlassToKey.Linux.Gui" "$@"' "${local_wrapper}" ||
     grep -Fq 'exec /opt/GlassToKey.Linux.Gui/GlassToKey.Linux.Gui "$@"' "${local_wrapper}" ||
     grep -Fq 'GlassToKey compatibility shim for legacy /usr/local PATH precedence.' "${local_wrapper}"; then
    rm -f "${local_wrapper}" || true
  fi
}

if ! getent group glasstokey >/dev/null 2>&1; then
  groupadd --system glasstokey || true
fi

if command -v udevadm >/dev/null 2>&1; then
  udevadm control --reload-rules || true
  udevadm trigger || true
fi

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload || true
fi

remove_legacy_local_gui_wrapper

cat <<'EOF'
GlassToKey Linux install complete.

Recommended next steps:
  1. Reconnect the trackpads or wait a few seconds for refreshed udev permissions.
  2. Add the desktop user to the 'glasstokey' group:
     sudo usermod -aG glasstokey $USER
  3. Log out and back in so the new group membership applies.
  4. Run 'glasstokey doctor'
  5. Run 'glasstokey init-config' if this is the first install
  6. Run 'glasstokey show-config --print' and verify left/right device bindings
  7. Run 'glasstokey start' to launch the background runtime
  8. Run 'glasstokey stop' to stop the background runtime
  9. Run 'glasstokey run-engine 10' for a direct foreground smoke test when needed

Optional user service:
  systemctl --user daemon-reload
  systemctl --user enable --now glasstokey.service

Optional GUI:
  glasstokey-gui            # start tray host in background
  glasstokey-gui --show     # open config window
EOF
